# systemd/codex-model-sync.timer
[Unit]
Description=Periodically synchronize Codex model catalog

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true
Unit=codex-model-sync.service

[Install]
WantedBy=timers.target
